# roles/frps/templates/frps.service.j2
[Unit]
Description=frps

[Service]
User=root
WorkingDirectory=/etc/frp
ExecStart=frps --config /etc/frp/frps.toml
Restart=always

[Install]
WantedBy=multi-user.target
