# roles/frpc/templates/frpc.service.j2
[Unit]
Description=frpc

[Service]
User={{ ansible_user }}
WorkingDirectory=/etc/frp
ExecStart=frpc \
  --config=frpc.toml
Restart=always

[Install]
WantedBy=multi-user.target
